# systemd/gv50-tracker-simple.service
[Unit]
Description=GV50 GPS Tracker Service - Simplified Structure
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=gv50tracker
Group=gv50tracker
WorkingDirectory=/tracker/gv50
Environment=PATH=/tracker/venv/bin
ExecStart=/tracker/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gv50-tracker

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tracker/logs
ReadWritePaths=/tracker/backup

[Install]
WantedBy=multi-user.target